# example/systemd_timer/test.service
[Unit]
Description=Test systemd timer

[Service]
Type=simple
ExecStart=%h/tmp/test.sh

[Install]
WantedBy=default.target
